# supervisor/salt.ini
[program:salt]
command={{ salt['config.get']('virtualenv_path', '/SaViEn') }}/bin/salt-minion
