# state-tree/supervisor/salt.ini
[program:salt]
command={{ salt['config.get']('virtualenv_path', '/SaViEn') }}/bin/salt-minion
autorestart=true
exitcodes=9
